# configs/systemd/user/swaync.service.d/override.conf
[Unit]
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
RestartSec=3
Restart=on-failure
KillMode=mixed